# .config/systemd/user/trash-empty@.service
[Unit]
Description=Clear trash older than %i days

[Service]
Type=oneshot
ExecStart=/usr/bin/trash-empty %i
