# wallpaper/wallpaper.timer
[Unit]
Description=Change the wallpaper daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
